# home/profile.d/common
# vim: ft=sh fdm=marker
# {{{ Aliases
# {{{ Debian
if [ -d /etc/apt ]; then
	alias apti="apt install --no-install-recommends";
	alias aptir="apt install";
	alias aptp="apt purge";
	alias aptr="apt remove";
    alias apts="apt search";
	alias aptss="aptitude search";
	alias aptu="apt update";
	alias aptup="apt upgrade";
fi
# }}}

# {{{ Archlinux
if [ -d /etc/pacman.d ]; then
	# Lists pacman packages sorted by size
	alias pac-list='expac -HM "%011m\t%-20n\t%10d" $( comm -23 <(pacman -Qqen|sort) <(pacman -Qqg base base-devel|sort) ) | sort -rn'

	# Pacman aliases
	alias paci='pacman -S'
	alias pacs='pacman -Ss'
    alias pacu='pacman -Syy'
	alias pacup='pacman -Syu'
	alias pacr='pacman -Rs'

    # Yay aliases
    alias yayi='yay -S'
    alias yays='yay -Ss'
    alias yayup='yay  -Syu --nodevel'
    alias yayup-devel='yay  -Syu --devel'

    # Octopi
    alias pacg='/usr/bin/octopi'
fi
# }}}

# {{{ Entware
if [ -f /opt/etc/opkg.conf  ]; then
    alias opks='/opt/bin/opkg list | grep -i'
    alias opki='opkg install'
    alias opkr='opkg remove'
    alias opku='opkg update'
    alias opkup='opkg upgrade'
fi
# }}}

# {{{ Systemd
# Its official now, systemd is everywhere so might aswell get used to it
if command -v "systemd" >/dev/null 2>&1; then
    alias sd='systemctl'
    alias sdu='systemctl --user'
fi
# }}}

# {{{ General aliases
alias grep='grep --color=auto'
alias p='ps aux|grep -i --color=auto'
alias sudo='sudo '
alias svim='sudo -e'
alias ssh="ssh -o StrictHostKeyChecking=no"
alias steam="steam-native -silent"
alias tmuxa="tmux attach || tmux"

# ls
alias sl="ls"
alias la='ls -a'
alias ll='ls -l'
alias ls='ls --group-directories-first --color=always -h'

# Stone Soup Crawl XDG Compliance
alias crawl='crawl -rc $XDG_CONFIG_HOME/crawl/config -rcdir $XDG_CONFIG_HOME/crawl/ -macro $XDG_CONFIG_HOME/crawl'
alias craw-tiles='crawl-tiles -rc $XDG_CONFIG_HOME/crawl/config -rcdir $XDG_CONFIG_HOME/crawl/ -macro $XDG_CONFIG_HOME/crawl'
# }}}
# }}}

# {{{ Function aliases
# Colored man page
man() {
    env LESS_TERMCAP_mb=$'\E[01;31m' \
    LESS_TERMCAP_md=$'\E[01;38;5;74m' \
    LESS_TERMCAP_me=$'\E[0m' \
    LESS_TERMCAP_se=$'\E[0m' \
    LESS_TERMCAP_so=$'\E[38;5;246m' \
    LESS_TERMCAP_ue=$'\E[0m' \
    LESS_TERMCAP_us=$'\E[04;38;5;146m' \
    man "$@"
}
# }}}
